# Debian build/NornsWFM/DEBIAN/prerm
#/bin/bash
echo Stopping NornsWFM Service
sudo systemctl disable NornsWFM.service
sudo systemctl stop NornsWFM.service
